# confsupervisor/ginco_supervisor_tpl.conf
; Supervisor config file for Ginco instance
; Manage RabbitMQ workers
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[program:ginco_@instance.name@]
command=php /var/www/@apache.path.to.app@/website/server/app/console rabbitmq:consumer -m 1 @consumer.name@
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/@apache.path.to.app@/website/server/app/logs/rabbitMQconsumers.log